# telegram-agent/kiezquiz-agent.service
[Unit]
Description=KiezQuiz Telegram Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jjl
WorkingDirectory=/home/jjl/KiezQuiz/telegram-agent
Environment=PATH=/home/jjl/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/jjl/KiezQuiz/telegram-agent/.venv/bin/python bot.py
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
